# scripts/agents/plists/com.macpilot.portfolio-meeting-prep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.macpilot.portfolio-meeting-prep</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MACPILOT_AGENT_NAME</key>
    <string>portfolio-meeting-prep</string>
    <key>PROJECT_DIR</key>
    <string>__HOME__/Projects/iteration_audit</string>
    <key>MEETING_DURATION</key>
    <string>30m</string>
    <key>MEETING_MODEL</key>
    <string>opus</string>
    <key>MEETING_MAX_TURNS</key>
    <string>20</string>
    <key>MEETING_TIMEOUT</key>
    <string>1200</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>__MACPILOT_DIR__/agents/portfolio-meeting-prep.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>50</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__MACPILOT_DIR__/logs/portfolio-meeting-prep.out</string>
  <key>StandardErrorPath</key>
  <string>__MACPILOT_DIR__/logs/portfolio-meeting-prep.err</string>
</dict>
</plist>
